# br2rockpis/board/rockpis/rootfs_overlay/etc/systemd/system/dnsmasq.service
[Unit]
Description=DNS cacher and DHCP server
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/dnsmasq -k
Restart=always

[Install]
WantedBy=multi-user.target
