# br2rockpis/board/rockpis/rootfs_overlay/etc/systemd/system/usbgadget.service
[Service]
type=oneshot
ExecStart=/opt/usb_gadget.sh
RemainAfterExit=true
StandardOutput=journal
